# amavisd-clean-tmp.timer
[Unit]
Description=Amavisd clean tmp timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=24h

[Install]
WantedBy=timers.target
